# SpinnakerUpdater.cron
*/1 * * * * ~/Desktop/InstallSpinnaker/SpinnakerUpdater.sh
