# SpinnakerUpdater.cron
0 8 * * * ~/SpinnakerUpdater/SpinnakerUpdater.sh
