# SpinnakerUpdater.cron
*/1 * * * * ~/SpinnakerUpdater/SpinnakerUpdater.sh
